# Server/userdbserver.service
[Unit]
Description=UserDbServer

[Service]
ExecStart=/opt/UserDbServer/bin/UserDbServer
Type=forking

[Install]
WantedBy=multi-user.target
