# Server/userdbserver.service
[Unit]
Description=UserDbServer

[Service]
ExecStart=/opt/Server/bin/Server
Type=forking

[Install]
WantedBy=multi-user.target
